# _ttp/com.rmgdri.security-check.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rmgdri.security-check.morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/rayrichardson/ControlHub/RMGDRI_Website/rmgdri-site/_ttp/run-automated-check.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/rayrichardson/ControlHub/RMGDRI_Website/rmgdri-site/_ttp/logs/launchd-morning-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rayrichardson/ControlHub/RMGDRI_Website/rmgdri-site/_ttp/logs/launchd-morning-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/rayrichardson/ControlHub/RMGDRI_Website/rmgdri-site</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>AUTO_COMMIT</key>
        <string>false</string>
    </dict>
</dict>
</plist>
